# bin/openrc-run
#!/sbin/openrc-run

name="isep"
description="Service to run Jobs4u application service"

# Use $HOME instead of ~ and : instead of ;
project="~/sem4pi-23-24-2di2"
command="./bin/run-server.sh"
pidfile="/var/run/isepService.pid"
command_background="yes"

depend() {
    need localmount
    use logger
    after bootmisc
}

start_pre() {
    # Ensure the PID file directory exists
    if [ ! -d "/var/run" ]; then
        mkdir -p /var/run
        chown root:root /var/run
        chmod 755 /var/run
    fi
}

start() {
    ebegin "Starting isep"
    start-stop-daemon --start --pidfile "${pidfile}" --exec /bin/sh -- -c "cd "${project}" && ${command}"
    eend $?
}

stop() {
    ebegin "Stopping isep"
    start-stop-daemon --stop --pidfile "${pidfile}" --name "java"
    eend $?
}